# systemd/hunter-memory.service
[Unit]
Description=Hunter Memory System Server
After=network.target

[Service]
Type=simple
User=REPLACE_WITH_USERNAME
WorkingDirectory=REPLACE_WITH_INSTALL_PATH
ExecStart=/usr/bin/python3 server.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
